# emulator/files/rm-qemu-emergency.service
[Unit]
ExecStart=/usr/bin/systemctl poweroff

[Install]
WantedBy=emergency.target
